# common/share/cips_resources/cipsschedop.cron
#-----------------------------------------------
# Log Rotation 
#-----------------------------------------------
58 23 * * * $COMMON_HOME/bin/taskcentre_logrotate.sh $SHARED_HOME/journal/sms/logs 7
